# sensor_project/crontab
@reboot sudo python sensor_box.py
@reboot sudo python force_reboot.py
